# package/base-files/files/etc/init.d/epoll
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=80
start()
{
sleep 5
/usr/sbin/epoll_test&
}
stop()
{
killall epoll_test
}
